# crontab
0 16 * * * bundle exec rake lottery:execute_lottery  
